# recipes-core/remotechargeport/remotechargeport/systemaggregatorftpd.socket
[Unit]
Description=FTP Daemon Socket for EVerest's SystemAggregator

[Socket]
ListenStream=2121
Accept=yes

[Install]
WantedBy=sockets.target
